# source/etc/systemd/system/systemd-networkd-wait-online.service.d/10-add-on-failure.conf
# When wait-online has failed, start the wifi-setup target
[Unit]
OnFailure=cfg-ap.target
